# contrib/systemd/veil.service
# ============================================
# veil systemd サービスファイル
# ============================================
#
# io_uring ベースの高性能リバースプロキシサーバー用
# サンドボックス化による堅牢なセキュリティを提供
#
# インストール:
#   sudo systemctl daemon-reload
#   sudo systemctl enable veil
#   sudo systemctl start veil
#
# ============================================

[Unit]
Description=Veil - High-Performance Reverse Proxy Server (io_uring + rustls)
Documentation=https://github.com/aofusa/veil
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/veil -c /var/etc/veil/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# ============================================
# ユーザー・グループ設定
# ============================================
User=veil
Group=veil

# ============================================
# リソース制限
# ============================================
LimitNOFILE=1048576
LimitMEMLOCK=infinity
LimitNPROC=4096
LimitCORE=0

# ============================================
# ケイパビリティ設定
# ============================================
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW
AmbientCapabilities=CAP_NET_BIND_SERVICE

# ============================================
# ファイルシステム隔離
# ============================================
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
LogsDirectory=veil
CacheDirectory=veil
ReadOnlyPaths=/var/etc/veil

# ============================================
# 名前空間隔離
# ============================================
RestrictNamespaces=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectProc=invisible

# ============================================
# ネットワーク設定
# ============================================
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

# ============================================
# セキュリティ強化
# ============================================
NoNewPrivileges=yes
# WASM (full feature) は実行可能メモリを必要とするため無効化
# MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes

# システムコール制限は config.toml の seccomp / Landlock に委譲する。
# systemd の SystemCallFilter と併用すると veil 側の seccomp/Landlock 適用が失敗する。

# ============================================
# 環境変数
# ============================================
Environment=RUST_BACKTRACE=1
Environment=MALLOC_CONF=background_thread:true

[Install]
WantedBy=multi-user.target