# packaging/bsd/openbsd/veil.rc
#!/bin/ksh
#
# OpenBSD rc.d サービススクリプト（veil、F-120 Phase 6）
#
# インストール:
#   install -m 0555 veil.rc /etc/rc.d/veil
#   rcctl enable veil
#   rcctl start veil
#
# 設定ファイル: /etc/veil/config.toml
#
# 注意: veil のアプリ内 pledge/unveil（[security] enable_pledge / enable_unveil）は
# veil 自身のプロセスを制限する。OpenBSD ネイティブのセキュリティを最大限使うには
# config.toml で enable_pledge = true / enable_unveil = true を設定すること。
# TLS は rustls の ring プロバイダで動作する（F-122）。HTTPS 静的配信/プロキシとも
# pledge+unveil 有効のまま 200 で動作することを検証済み。

daemon="/usr/local/bin/veil"
daemon_flags="-c /etc/veil/config.toml"
daemon_user="root"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_reload() {
    pkill -HUP -xf "${daemon} ${daemon_flags}"
}

rc_check() {
    pgrep -xf "${daemon} ${daemon_flags}" >/dev/null
}

rc_cmd $1
